# router.service
[Unit]
Description=Configure router

[Service]
Type=simple
ExecStart=/opt/network

[Install]
WantedBy=multi-user.target